# usr.xtc/xtclogin/xtclogin.service
[Service]
Name=XTC Window Manager
ExecStart=/usr/xtc/bin/xtclogin
RunLevel=5
