# etc/doit.d/devconsole.service
[Service]
Name=console
ExecStart=/sbin/newtty /dev/ttyS0
RunLevel=2
