# etc/doit.d/devconsole.service
[Service]
Name=devconsole
ExecStart=/sbin/newtty /dev/ttyS1
