# etc/doit.d/systerm.service
[Service]
Name=systerm
ExecStart=/usr/libexec/systerm
RunLevel=2
