# etc/doit.d/systerm.service
[Service]
Name=systerm
ExecStart=/sbin/systerm
