# etc/doit.d/systerm.multiuser.service
[Service]
Name=systerm.multiuser
ExecStart=/usr/libexec/systerm
RunLevel=2
