# etc/doit.d/systerm.singleuser.service
[Service]
Name=systerm.singleuser
ExecStart=/usr/libexec/systerm /bin/sh
RunLevel=1
