# etc/doit.d/xtcwm.service
[Service]
Name=XTC Window Manager
ExecStart=/bin/sh /etc/xtc/xtcwm.rc
RunLevel=5
