# etc/doit.d/xtcwm.service
[Service]
Name=XTC Window Manager
ExecStart=/usr/libexec/xtcwm
RunLevel=5
